# etc/supervisord/supervisor.conf
[inet_http_server]
port = 0.0.0.0:9001

[supervisord]
logfile = /var/log/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 1
loglevel = info
pidfile = /var/run/supervisord.pid
nodaemon = false
minfds = 800
childlogdir = /var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]

[program:elasticsearch]
command = /usr/share/elasticsearch/bin/elasticsearch
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stopasgroup = true
umask=022
priority=10
autostart=true
autorestart=true
startsecs=8
startretries=3
stopsignal=KILL
stopwaitsecs=5
user=elasticsearch
redirect_stderr=true


[program:kibana]
command = /opt/kibana/bin/kibana
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stopasgroup = true
umask=022
priority=30
autostart=true
autorestart=true
startsecs=8
startretries=3
stopsignal=KILL
stopwaitsecs=5
# user=kibana
redirect_stderr=true


[program:logstash]
command = /opt/logstash/bin/logstash agent -f /etc/logstash/conf.d
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stopasgroup = true
umask=022
priority=20
autostart=true
autorestart=true
startsecs=8
startretries=3
stopsignal=KILL
stopwaitsecs=5
user=logstash
redirect_stderr=true
